# Library/LaunchAgents/com.dotfiles.chezmoi-sync.plist.tmpl
{{- if not (hasKey . "autosync") -}}{{ fail "missing autosync template var; set [data].autosync = true/false in chezmoi config" }}{{- end -}}
{{- if .autosync -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.dotfiles.chezmoi-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>{{ .chezmoi.homeDir }}/.local/bin/chezmoi-sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/chezmoi-sync.launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/chezmoi-sync.launchd.err.log</string>
  </dict>
</plist>
{{- end -}}
